# image/files/etc/profile.d/20-google-cloud-sdk.sh
pathmunge /opt/google-cloud-sdk/bin after
